# specs/default/chef/site-cookbooks/anaconda/templates/default/jupyter.service.erb
# After Ubuntu 16.04, Systemd becomes the default. 
# It is simpler than https://gist.github.com/Doowon/38910829898a6624ce4ed554f082c4dd

[Unit]
Description=Jupyter Notebook

[Service]
Type=simple
PIDFile=/run/jupyter.pid
ExecStart=<%= node['anaconda']['install_root'] %>/<%= node['anaconda']['version'] %>/bin/jupyter-notebook  --no-browser --config=<%= node['anaconda']['home'] %>/.jupyter/jupyter_notebook_config.py --port=<%= node['anaconda']['notebook']['port'] %> --notebook-dir=<%= node['anaconda']['home'] %>/.local
User=<%= node['anaconda']['owner'] %>
Group=<%= node['anaconda']['owner'] %>
WorkingDirectory=<%= node['anaconda']['home'] %>
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=jupyter
#KillMode=mixed

[Install]
WantedBy=multi-user.target